# includes/meshcentral
#!/bin/sh

# MeshCentral FreeBSD Service Script

# PROVIDE: meshcentral
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr                                                                                       

load_rc_config $name                                                                                 

: ${meshcentral_enable:="NO"}
: ${meshcentral_pidfile:="/var/run/meshcentral.pid"}
: ${meshcentral_chdir:="/usr/local/node_modules"}                                                                                              

name="meshcentral"                                                                                   
rcvar=meshcentral_enable                                                                             
command="/usr/sbin/daemon"                                                                           
PATH=${PATH}:/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin:/root/bin                 
procname="/usr/local/bin/node"
command_args="-r -P ${meshcentral_pidfile} ${procname} ${meshcentral_chdir}/${name}"
                                                                                                    
run_rc_command "$1"
